# systemd-templates/.config/systemd/user/dotfiles-backup.service
[Unit]
Description=Backup dotfiles and Obsidian

[Service]
Type=oneshot
ExecStart=__HOME_DIR__/.local/bin/backup.sh
