# systemd-templates/.config/systemd/user/sync-myshare.service
[Unit]
Description=Sync MyShare folder to the cloud

[Service]
Type=oneshot
ExecStart=__HOME_DIR__/.local/bin/sync-myshare.sh
