# systemd-templates/.config/systemd/user/sync-myshare.timer
[Unit]
Description=Run MyShare sync every 4 hours

[Timer]
OnUnitInactiveSec=4h
Persistent=true

[Install]
WantedBy=timers.target
