# systemd-system-drop-ins/ollama.service.d/10-cors-override.conf
# /home/imak/dotfiles/systemd-system-drop-ins/ollama.service.d/10-cors-override.conf
[Service]
Environment="OLLAMA_ORIGINS=*"
